# systemd/.config/systemd/user/loq-cliphist-image.service
[Unit]
Description=Cliphist image watcher
PartOf=loq-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/sh -lc 'wl-paste --type image --watch cliphist store'
Restart=always
RestartSec=1

[Install]
WantedBy=loq-session.target

