# systemd/.config/systemd/user/lumina-shell.service
[Unit]
Description=Lumina Shell
After=graphical-session.target
PartOf=loq-session.target
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.local/bin/lumina-shell
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

[Install]
WantedBy=loq-session.target
