# systemd/.config/systemd/user/lumina-wallpaper-rotate.timer
[Unit]
Description=Check the Lumina wallpaper rotation schedule

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=5s
Unit=lumina-wallpaper-rotate.service

[Install]
WantedBy=timers.target
